# src-tauri/embeds/macos/com.clamav-desktop.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clamav-desktop.daemon</string>

  <key>ProgramArguments</key>
  <array>
      <string>/usr/local/bin/clamav-desktop-daemon</string>
  </array>

  <key>UserName</key>
  <string>root</string>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/clamav-desktop-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/clamav-desktop-daemon.err</string>
</dict>
</plist>
